# backend/deploy/systemd/backend-api.service
[Unit]
Description=AutoSites backend API (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/AutoSites/backend
EnvironmentFile=/home/%i/AutoSites/backend/.env
ExecStart=/home/%i/.local/bin/uv run uvicorn api.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
